# ext/profile.d/UMobj.csh
if ( $uid != 0 ) then
    modulecmd csh add UMobj
endif
